# etc/supervisor/conf.d/terminal.conf
[program:terminal]
command=x-terminal-emulator -geometry 120x40+10+10 -ls -title "x-terminal-emulator"
priority=15
redirect_stderr=true
stdout_logfile_backups=2
stdout_logfile_maxbytes=10MB
