# etc/supervisor/conf.d/vnc.conf
[program:vnc]
command=x11vnc -forever -usepw -display %(ENV_DISPLAY)s
priority=15
redirect_stderr=true
stdout_logfile_backups=2
stdout_logfile_maxbytes=10MB
